# scripts/backup/bin/basex_backup.cron
#use rsync to aggregate backup files
MAILTO="jeffery.antoniuk@ualberta.ca"
33 03 * * * basex /var/basex/backup/bin/backup_basex.sh
59 03 * * * basex /var/basex/backup/bin/backup_basex_rotate.sh

